# setup/local-git-daemon.conf
start on startup
stop on shutdown
exec /usr/bin/git daemon --base-path=/repos --export-all --enable=receive-pack --reuseaddr --informative-errors --verbose --detach
respawn
